# docker/focal/supervisor/conf.d/mysql.conf
[program:mysql]
# /usr/bin/mysqld_safe
command=/usr/sbin/mysqld
autostart=true
user=root
autorestart=true
loglevel=info
