# docker/focal/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
user=root
autorestart=true
loglevel=info
