# docker/ubuntu/supervisor/minio.conf
[program:minio]
command="MINIO_ROOT_USER=admin MINIO_ROOT_PASSWORD=12345678 /home/deploy/.local/bin/minio server /var/lib/minio/data"
autostart=true
user=deploy
autorestart=true
loglevel=info

