# docker/ubuntu/supervisor/mysql.conf
[program:mysql]
command=/usr/sbin/mysqld
# command=mysqld_safe --skip-grant-tables
autostart=true
user=mysql
autorestart=true
loglevel=info
