# docker/ubuntu/supervisor/nginx.conf
[program:php-fpm]
command=/usr/sbin/nginx
autostart=true
user=root
autorestart=true
loglevel=info
