# docker/ubuntu/supervisor/rabbitmq.conf
[program:postgresql]
directory=/var/lib/rabbitmq
command=/usr/lib/rabbitmq/bin/rabbitmq-server
autostart=true
user=rabbitmq
group=rabbitmq
autorestart=true
loglevel=info
