# docker/ubuntu/supervisor/redis.conf
[program:redis]
directory=/var/lib/redis
command=/usr/bin/redis-server /etc/redis/redis.conf
autostart=true
user=root
autorestart=true
loglevel=info
