# docker/ubuntu/supervisor/conf.d/elastic-search.conf
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
autostart=true
user=elasticsearch
autorestart=true
loglevel=info
