# docker/ubuntu/supervisor/conf.d/mysql.conf
[program:mysql]
command=/usr/sbin/mysqld
# command=mysqld_safe --skip-grant-tables
autostart=true
user=root
autorestart=true
loglevel=info
