# docker/ubuntu/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm8.0 -F
autostart=true
user=root
autorestart=true
loglevel=info

